# buildroot/fsroot/etc/rc.local
#!/bin/sh

export PATH=/bin:/sbin:/usr/bin:/usr/sbin

true
exit 0
